# tools/analytics_server/ansible/roles/server/templates/server.ini.j2
[program:analytics_server]
command="{{ exec_path }}/analytics_server.py"
autostart=true
autorestart=true
stopsignal=KILL
stopasgroup=true
killasgroup=true

